# build/sprvsr_uwc1.conf
[program:uwc_stage1]
command=/vagrant/env/bin/uwsgi /vagrant/wsgi.ini
autostart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=INT

[program:uwc_stage1_rq1]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq2]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq3]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq4]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq5]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq6]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq7]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq8]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq9]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT

[program:uwc_stage1_rq10]
command=/vagrant/env/bin/python /vagrant/manage.py rqworker high default low
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=QUIT
